# pkgbuild/anti-ruscher.install
post_install() {
    systemctl daemon-reload
    systemctl enable anti-ruscher.service
}

pre_remove() {
    systemctl disable anti-ruscher.service
} 
