# pkgbuild/anti-ruscher.install
# post_install() {
#         systemctl daemon-reload
#         systemctl enable --global anti-ruscher.service
# }
# 
# pre_remove() {
#     systemctl disable anti-ruscher.service
# } 
