# usr/lib/systemd/user/anti-ruscher.timer
 
[Unit]
Description=Anti-Ruscher

[Timer]
OnUnitActiveSec=120min

[Install]
WantedBy=multi-user.target 
